# config/etc/cron.d/wikilegis
30 5 * * * python3 /var/labhacker/wikilegis/wikilegis/manage.py newsletter_notification --periodicity daily
0 6 * * 6 python3 /var/labhacker/wikilegis/wikilegis/manage.py newsletter_notification --periodicity weekly
